# docker/conf.d/afs-celerybeat.conf
; ================================
;  celery beat supervisor
; ================================

[program:celerybeat]
    command=python3 -m celery -A arches_for_science.celery beat --loglevel=INFO 
    directory=/web_root/afs

    user=root
    numprocs=1
    stdout_logfile=/var/log/celery/beat.log
    stderr_logfile=/var/log/celery/beat.log
    autostart=true
    autorestart=true
    startsecs=10

    ; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
    stopasgroup=true

    ; if rabbitmq is supervised, set its priority higher
    ; so it starts first
    priority=999
